# build/teleshell.service
[Unit]
Description=TeleShell Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/opt/bin/teleshell
Environment="TELESHELL_API_TOKEN=*********"
Environment="TELESHELL_PASSWORD=**********"
Environment="TELESHELL_BASH_PATH=/bin/bash"

[Install]
WantedBy=multi-user.target
